# roles/master/files/kube-scheduler.service
[Unit]
Description=kube-controller-manager
After=kube-apiserver.service

[Service]
EnvironmentFile=-/etc/default/%p
ExecStart=/usr/bin/kube-scheduler $KUBE_SCHEDULER_OPTS
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
